# service/dbus-tasmota-inverter/run
#!/bin/sh
exec 2>&1
exec /data/dbus-tasmota-inverter/dbus-tasmota-inverter.py
